# vmware-tools/vmware-tools-distrib/installer/upstart-job.conf
##########################################################################
# Copyright (c) 2010-2015 VMware, Inc.  All rights reserved.
##########################################################################

##########################################################################
# DO NOT modify this file directly as it will be overwritten the next
# time the VMware Tools are installed.
##########################################################################

description     "VMware Tools services"
author          "VMware, Inc."

# Be sure to block the display managers until our job has completed.  This
# is to make sure our kernel services are running before vmware-user
# may launch.
start on runlevel [235] or starting gdm or starting kdm or starting prefdm
stop on runlevel [06]

pre-start exec /etc/vmware-tools/services.sh start
post-stop exec /etc/vmware-tools/services.sh stop
